# deploy/systemd/aifar-runtime.service
[Unit]
Description=AIFAR Runtime
Wants=network-online.target
Requires=docker.service
After=network-online.target docker.service

[Service]
Type=simple
User=aifar-runtime
Group=aifar-runtime
SupplementaryGroups=docker
UMask=0027
WorkingDirectory=/
Environment=AIFAR_RUNTIME_CONFIG=/etc/aifar-runtime/config.yaml
EnvironmentFile=-/etc/aifar-runtime/aifar-runtime.env
ExecStart=/usr/local/bin/aifar-runtime serve --config ${AIFAR_RUNTIME_CONFIG}
Restart=on-failure
RestartSec=3s
KillSignal=SIGTERM
TimeoutStopSec=35s
LimitNOFILE=1048576
StateDirectory=aifar-runtime
LogsDirectory=aifar-runtime
ConfigurationDirectory=aifar-runtime
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/aifar-runtime /var/log/aifar-runtime
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
